# root/etc/services.d/lidarr/run
#!/usr/bin/with-contenv bash

umask 022

cd /opt/Lidarr || exit

exec \
	s6-setuidgid appuser mono Lidarr.exe \
	-nobrowser -data=/config
